# agent/systemd/led-vfx-agent.service
[Unit]
Description=LED VFX creativity agent (one session)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ejw
# Plain paths, not %h - %h resolves relative to the SERVICE MANAGER's own
# context (effectively root) on a system unit like this one, not to the
# User= set here. Confirmed the hard way: relying on it produced "Failed
# to load environment files" / "Failed to spawn 'start' task: No such
# file or directory" (systemd was looking under /root, not /home/ejw).
WorkingDirectory=/home/ejw/Code/led-vfx
EnvironmentFile=/home/ejw/Code/led-vfx/.env
ExecStart=/usr/local/bin/node agent/session.js
# Generous but bounded - a hung API call/stream shouldn't run forever.
TimeoutStartSec=1800
# Deliberately no Restart= : a failure (exhausted retries, bad key, network
# blip) should wait for the next hourly timer trigger, not retry in a tight
# loop and compound API cost on a persistent problem. A failed run shows up
# in `systemctl --failed` / `journalctl -u led-vfx-agent.service` for
# whoever checks in on it - that's the intended visibility, not silence.
